# config/supervisord.conf
;/*
; * This file is part of kgiflwl/v2board-compose.
; *
; * kgiflwl/v2board-compose is free software: you can redistribute it and/or modify
; * it under the terms of the GNU General Public License as published by
; * the Free Software Foundation, either version 3 of the License, or
; * (at your option) any later version.
; *
; * kgiflwl/v2board-compose is distributed in the hope that it will be useful,
; * but WITHOUT ANY WARRANTY; without even the implied warranty of
; * MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the
; * GNU General Public License for more details.
; *
; * You should have received a copy of the GNU General Public License
; * along with kgiflwl/v2board-compose. If not, see <https://www.gnu.org/licenses/>.
; */

;supervisor config

[program:v2board]
process_name=%(program_name)s_%(process_num)02d
command=php /www/artisan horizon
stdout_logfile=/wwwlogs/horizon.log
stdout_logfile_maxbytes=0
stderr_logfile=/wwwlogs/horizon_error.log
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startretries=0
numprocs=2
